# deploy/tunnel-watchdog@.timer
[Unit]
Description=Run tunnel watchdog every 5 minutes for %i
After=gateway@%i.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
